# autostart/jamlinux-startup-sound.desktop
[Desktop Entry]
Type=Application
Name=JamLinux Greeter Startup Sound
Comment=Play the branded startup sound when the login screen is ready
TryExec=/usr/local/bin/jamlinux-startup-sound
Exec=/usr/local/bin/jamlinux-startup-sound
Terminal=false
NoDisplay=true
X-GNOME-Autostart-Delay=5
