# systemd/jamlinux-first-boot.service
[Unit]
Description=Run JamLinux first-boot tasks
Wants=network-online.target
After=network-online.target
ConditionPathExists=!/var/lib/jamlinux/first-boot-complete

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/first-boot.sh
Restart=on-failure
RestartSec=2min
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
